# config/rmapp.ini
[group:rmappgroup]
programs=rmapp

[program:rmapp]
directory=/home/mitch/app/rm                     ; Location of application
command=/home/mitch/app/rm/env/bin/python3.9 main.py                   ; The command to execute
autostart=true                           ; Start this application when supervisord starts
autorestart=true                         ; Restart this application if it crashes
stderr_logfile=/var/log/rmapp/app.err.log  ; Make sure this directory exists
stdout_logfile=/var/log/rmapp/app.log      ; Make sure this directory exists
stopsignal=INT                           ; Signal sent to the application when halting
user=mitch                           ; setuid to this UNIX account to run the program
stopwaitsecs=610                        ; Wait 10 Minutes for the threads to exit and main to return.